# packs/mcp/mcp-persistence.yaml
name: "MCP Persistence Evasion"
description: "Block MCP write_file tool calls to shell startup scripts and user-level persistence paths — LaunchAgents, systemd user services, and shell dotfiles"
version: "1.0.0"
author: "AgentShield Security Research"

# Shell-level rules block `echo >> ~/.bashrc`, `cat > ~/.zshrc`, and similar
# commands. But an AI agent using an MCP write_file tool bypasses shell
# evaluation entirely — there is no shell to intercept.
#
# These rules close the gap by blocking direct MCP writes to home-directory
# persistence paths that are NOT covered by the system-directory write rules
# in mcp-safety.yaml (/etc/**, /var/**).
#
# MITRE techniques:
#   T1546.004 — Event Triggered Execution: Unix Shell Configuration Modification
#   T1543.001 — Create or Modify System Process: Launch Agent
#   T1543.002 — Create or Modify System Process: Systemd Service

structural_rules:
  # === BLOCK: Shell startup dotfile modification via MCP write_file ===
  #
  # Shell rc/profile files execute on every new shell session. An agent that
  # writes malicious code here achieves persistent execution without needing
  # a cron job or systemd service. This is the simplest and most common
  # persistence technique on Unix/macOS systems.
  #
  # Covered files (home-directory dotfiles only — /etc/ paths are blocked by
  # mcp-safety-block-etc-write):
  #   Bash:  .bashrc, .bash_profile, .bash_login, .bash_logout
  #          .bash_aliases (sourced by default Ubuntu/Debian .bashrc)
  #          .bash_functions (commonly sourced by .bashrc)
  #   Zsh:   .zshrc, .zprofile, .zlogin, .zlogout, .zshenv
  #          .zsh_aliases (commonly sourced in zsh configs)
  #          .zshenv — sourced for ALL zsh processes including non-interactive
  #   POSIX: .profile
  #   Other: .kshrc, .cshrc, .tcshrc
  - id: mcp-persist-block-shell-profile-write
    taxonomy: "persistence-evasion/scheduled-tasks/crontab-modification"
    match:
      tool_name_any:
        - "write_file"
        - "create_file"
        - "edit_file"
        - "write_to_file"
        - "save_file"
        - "update_file"
        - "append_file"
        - "filesystem_write"
      args_match:
        path:
          pattern_any:
            - "(?:^|/)\\.(bashrc|zshrc|bash_profile|profile|bash_login|bash_logout|zprofile|zlogin|zlogout|zshenv|kshrc|cshrc|tcshrc|bash_aliases|bash_functions|zsh_aliases)$"
    decision: "BLOCK"
    reason: "MCP write to shell startup dotfile — code injected here executes on every new shell session. Persistent execution without any cron/systemd entry. MITRE T1546.004."
    confidence: 0.95

  # === BLOCK: macOS LaunchAgent / LaunchDaemon persistence via MCP ===
  #
  # Writing a plist file to ~/Library/LaunchAgents/ registers a persistent
  # background service that launchd starts on user login. Writing to
  # /Library/LaunchDaemons/ (system-wide, root-owned) is already blocked by
  # mcp-safety-block-etc-write via /usr/ and by this rule as belt-and-suspenders.
  #
  # This is the canonical macOS persistence technique (T1543.001).
  - id: mcp-persist-block-launchagent-write
    taxonomy: "persistence-evasion/service-management/systemctl-manipulation"
    match:
      tool_name_any:
        - "write_file"
        - "create_file"
        - "edit_file"
        - "write_to_file"
        - "save_file"
        - "update_file"
        - "filesystem_write"
      args_match:
        path:
          pattern_any:
            - "Library/LaunchAgents/"
            - "Library/LaunchDaemons/"
    decision: "BLOCK"
    reason: "MCP write to macOS LaunchAgent/Daemon directory installs a persistent background service loaded on login. MITRE T1543.001."
    confidence: 0.95

  # === BLOCK: Git hook injection via MCP write_file (issue #150) ===
  #
  # Git hooks are executable scripts stored in .git/hooks/. When an agent
  # writes a malicious payload to pre-commit, post-checkout, post-merge, or
  # any other hook, the payload executes on every corresponding git operation
  # performed by the developer — without any shell command that AgentShield
  # can intercept (the write_file call never goes through a shell).
  #
  # This is a persistent code execution technique that:
  #   - Requires no elevated privileges (hooks run as the current user)
  #   - Survives process restarts (it's a file on disk)
  #   - Fires silently inside developer workflows (git commit, git pull, etc.)
  #   - Has no shell-level detection gap — only MCP can catch the write
  #
  # Standard git hook names:
  #   pre-commit, prepare-commit-msg, commit-msg, post-commit
  #   pre-push, post-push, pre-rebase
  #   post-checkout, post-merge, post-receive
  #   pre-receive, update, post-update
  #   applypatch-msg, pre-applypatch, post-applypatch
  #   sendemail-validate, fsmonitor-watchman
  #
  # MITRE T1546 — Event Triggered Execution (hooks fire on git operations)
  # MITRE T1059.004 — Command and Scripting Interpreter: Unix Shell
  - id: mcp-persist-block-git-hook-write
    taxonomy: "persistence-evasion/cicd-tampering/pipeline-config-write"
    match:
      tool_name_any:
        - "write_file"
        - "create_file"
        - "edit_file"
        - "write_to_file"
        - "save_file"
        - "update_file"
        - "append_file"
        - "filesystem_write"
        - "str_replace_editor"
        - "str_replace_based_edit_tool"
      args_match:
        path:
          pattern_any:
            - "(?:^|/)\\.git/hooks/"
    decision: "BLOCK"
    reason: "MCP write to .git/hooks/ installs a persistent git hook that executes on every matching git operation (commit, push, checkout, etc.) without any shell-level interception. MITRE T1546."
    confidence: 0.97

  # === BLOCK: Linux user-level systemd service and autostart persistence ===
  #
  # ~/.config/systemd/user/ holds user-level systemd unit files. A service
  # installed here starts automatically without root privileges via
  # `systemctl --user enable`. ~/.config/autostart/ holds XDG autostart
  # entries that desktop environments launch at login.
  #
  # These paths are not covered by mcp-safety-block-etc-write (/etc/systemd/
  # is covered there) because user-level services live in the home directory.
  - id: mcp-persist-block-user-systemd-write
    taxonomy: "persistence-evasion/service-management/systemctl-manipulation"
    match:
      tool_name_any:
        - "write_file"
        - "create_file"
        - "edit_file"
        - "write_to_file"
        - "save_file"
        - "update_file"
        - "filesystem_write"
      args_match:
        path:
          pattern_any:
            - "\\.config/systemd/user/"
            - "\\.config/autostart/"
            - "\\.local/share/systemd/user/"
    decision: "BLOCK"
    reason: "MCP write to user-level systemd/autostart directory installs a persistent service that starts without root on next login. MITRE T1543.002."
    confidence: 0.95

  # === BLOCK: SSH authorized_keys write via MCP (issue #151) ===
  #
  # Writing a public key to ~/.ssh/authorized_keys grants permanent remote
  # shell access to whoever holds the matching private key — no password
  # required. Writing to ~/.ssh/config can redirect all SSH connections
  # through an attacker-controlled proxy, enabling credential interception.
  #
  # Note: mcp-safety.yaml already broadly blocks all .ssh/** writes via the
  # mcp-safety-block-ssh-read rule. This rule adds specific taxonomy, reason
  # messages, and confidence scores for the highest-impact sub-paths, ensuring
  # clear audit records and rule attribution for these targeted attack vectors.
  # An MCP write_file call bypasses the shell entirely.
  #
  # Covered paths:
  #   ~/.ssh/authorized_keys — adds a persistent backdoor public key
  #   ~/.ssh/config          — can redirect or proxy SSH connections
  #   ~/.ssh/known_hosts     — can be poisoned to suppress MITM warnings
  #   /home/*/.ssh/authorized_keys — any user's authorized_keys
  #   /root/.ssh/authorized_keys   — root backdoor (requires root-level MCP)
  #
  # MITRE T1098.004 — Account Manipulation: SSH Authorized Keys
  # MITRE T1071.002 — Application Layer Protocol: File Transfer Protocols
  - id: mcp-persist-block-ssh-authkeys-write
    taxonomy: "persistence-evasion/lateral-movement/ssh-config-tampering"
    match:
      tool_name_any:
        - "write_file"
        - "create_file"
        - "edit_file"
        - "write_to_file"
        - "save_file"
        - "update_file"
        - "append_file"
        - "filesystem_write"
        - "str_replace_editor"
        - "str_replace_based_edit_tool"
      args_match:
        path:
          pattern_any:
            - "(?:^|/)\\.ssh/authorized_keys$"
            - "(?:^|/)\\.ssh/authorized_keys2$"
            - "(?:^|/)\\.ssh/config$"
            - "(?:^|/)\\.ssh/known_hosts$"
    decision: "BLOCK"
    reason: "MCP write to SSH authorized_keys adds a persistent backdoor public key granting remote shell access with no password. Write to ~/.ssh/config can proxy all SSH connections to an attacker. MITRE T1098.004."
    confidence: 0.97

  # === AUDIT: AI agent instruction and memory file writes via MCP (issues #206, #603) ===
  #
  # AI agent instruction files govern agent behavior across sessions. An MCP
  # write_file call targeting these files can inject persistent instructions
  # that alter agent behavior in all future sessions — a form of inter-session
  # prompt injection that bypasses all shell-level defenses.
  #
  # Covered files (matches ts-block-agent-instruction-inject in terminal-safety.yaml):
  #   Claude Code:     CLAUDE.md, .claude/CLAUDE.md, .claude/memory/*.md
  #   Cursor:          .cursorrules, .cursor/rules/
  #   Windsurf:        .windsurfrules
  #   Gemini CLI:      GEMINI.md
  #   OpenAI Codex:    AGENTS.md
  #   Amp:             AMP.md
  #   Cline:           .clinerules
  #   Roo Code:        .roomodes, .roo/
  #   GitHub Copilot:  .github/copilot-instructions.md
  #   Continue.dev:    .continue/config.json
  #   Aider:           .aider.conf.yml
  #
  # Decision: AUDIT (not BLOCK) — agents legitimately update CLAUDE.md and
  # memory files as part of normal operation (e.g., remembering user preferences).
  # Human review ensures no malicious instructions were injected.
  #
  # OWASP LLM01 — Prompt Injection
  # OWASP LLM08 — Integrity Failures
  # MITRE T1564 — Hide Artifacts (instructions hidden in memory files)
  - id: mcp-persist-audit-instruction-file-write
    taxonomy: "persistence-evasion/agent-config-persistence/agent-instruction-file-tampering"
    match:
      tool_name_any:
        - "write_file"
        - "create_file"
        - "edit_file"
        - "write_to_file"
        - "save_file"
        - "update_file"
        - "append_file"
        - "filesystem_write"
        - "str_replace_editor"
        - "str_replace_based_edit_tool"
      args_match:
        path:
          pattern_any:
            - "(?:^|/)CLAUDE\\.md$"
            - "(?:^|/)GEMINI\\.md$"
            - "(?:^|/)AGENTS\\.md$"
            - "(?:^|/)AMP\\.md$"
            - "(?:^|/)\\.cursorrules$"
            - "(?:^|/)\\.cursor/rules/"
            - "(?:^|/)\\.windsurfrules$"
            - "(?:^|/)\\.clinerules$"
            - "(?:^|/)\\.roomodes$"
            - "(?:^|/)\\.roo/"
            - "(?:^|/)\\.github/copilot-instructions\\.md$"
            - "(?:^|/)\\.continue/config\\.json$"
            - "(?:^|/)\\.aider\\.conf\\.yml$"
            - "(?:^|/)\\.claude/memory/.*\\.md$"
            - "(?:^|/)\\.claude/CLAUDE\\.md$"
    decision: "AUDIT"
    reason: "MCP write to AI agent instruction or memory file — injecting malicious instructions poisons future agent sessions (persistent inter-session prompt injection). Covers Claude Code, Cursor, Windsurf, Gemini CLI, Codex, Cline, Roo, Copilot, Continue.dev, Aider. OWASP LLM01/LLM08."
    confidence: 0.85
