# taxonomy/reconnaissance/system-enumeration/cron-config-read.yaml
id: "reconnaissance/system-enumeration/cron-config-read"
version: "1.0.0"
kingdom: "Reconnaissance & Discovery"
kingdom_id: 8
category: "System Enumeration"
category_id: "8.9"
name: "Cron Configuration Reconnaissance"
risk_level: "medium"

abstract: |
  An AI agent reads cron configuration files (/etc/crontab, /etc/cron.d/,
  /var/spool/cron/) to discover scheduled tasks, their execution schedules,
  and the scripts/commands they invoke — enabling privilege escalation via
  writable script discovery, evasion via monitoring schedule awareness, and
  persistence planning via injection into existing scheduled jobs.

explanation: |
  Cron configuration files define the scheduled execution landscape of a
  system. Reading these files gives an attacker or rogue AI agent three
  distinct attack advantages:

  **1. Privilege Escalation via Script Injection**
  Cron jobs often run as root or privileged service accounts. If a scheduled
  script is writable (misconfigured permissions), an agent that identifies
  the script via cron config can inject malicious code that executes as root
  in the next scheduled run. Key targets:
  - `/etc/cron.d/` — drop-in cron files, often with looser permissions
  - Scripts under `/etc/cron.{hourly,daily,weekly,monthly}/` — sometimes
    world-writable in misconfigured systems
  - `/var/spool/cron/root` — root's personal crontab

  **2. Evasion via Monitoring Schedule Discovery**
  Security tools often run on cron schedules (antivirus scans, log analysis,
  integrity checks). Discovering when `aide`, `rkhunter`, `lynis`, or custom
  security scripts run allows an attacker to time operations to avoid detection.

  **3. Persistence Reconnaissance**
  Understanding existing cron jobs helps identify:
  - Which jobs could be hijacked for persistence without adding new cron entries
  - Whether there are backup or cleanup jobs that might remove attacker artifacts
  - The overall system administration cadence (backup frequency, update schedules)

  **Common cron locations**:
  - `/etc/crontab` — system-wide crontab with user column
  - `/etc/cron.d/` — package-installed cron jobs
  - `/etc/cron.{hourly,daily,weekly,monthly}/` — run-parts directories
  - `/var/spool/cron/` (Linux) or `/var/spool/cron/crontabs/` — per-user crontabs

recommendation: |
  - **AUDIT** all AI agent reads of cron configuration paths.
  - Legitimate automation tasks requiring schedule awareness should be
    implemented through documented APIs, not by reading raw cron files.
  - Ensure scripts invoked by cron jobs as root are not world-writable;
    use a configuration management audit to catch permission drift.
  - Monitor cron config reads followed by attempts to write to script files
    (stateful: cron read → script write → privilege escalation chain).

examples:
  bad:
    - |
      # Discover all scheduled tasks and their scripts
      cat /etc/crontab && ls /etc/cron.d/ && ls /etc/cron.daily/
    - |
      # Check if any cron scripts are writable
      cat /etc/cron.d/backup-script
    - |
      # Enumerate root crontab
      cat /var/spool/cron/root
  good:
    - |
      # Create an application-specific cron entry (via crontab -e, not reading others)
      crontab -e

compliance:
  owasp-llm-2025: ["LLM06"]
  iso-27001-2022: ["A.8.16", "A.8.12"]
  soc2-type2: ["CC7.2", "CC6.1"]
  iso-42001-2023: ["A.9.4"]

references:
  mitre_attack: ["T1053.003", "T1057", "T1087"]
  cwe: ["CWE-200", "CWE-732"]

analyzers: ["mcp"]
related_rules: ["mcp-recon-audit-cron-config-read"]
